# wittleguys-newsletter/wittleguys-newsletter.service
[Unit]
Description=WittleGuys Newsletter and Contact Form Backend
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/wittleguys-newsletter
Environment=PATH=/opt/wittleguys-newsletter/venv/bin
ExecStart=/opt/wittleguys-newsletter/venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5000 app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

# Environment variables (set these in /etc/environment or use EnvironmentFile)
Environment=SMTP_SERVER=your-smtp-server.com
Environment=SMTP_PORT=587
Environment=SMTP_USER=your-smtp-username
Environment=SMTP_PASS=your-smtp-password
Environment=SECRET_KEY=your-secret-key-change-this
Environment=ADMIN_TOKEN=your-admin-token-change-this

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/wittleguys-newsletter

[Install]
WantedBy=multi-user.target
